# ansible/roles/fire2_worker/templates/fire2-worker@.service.j2
[Unit]
Description=fire2 worker node (instance %i)
Documentation=https://github.com/nallanos/fire2
After=network.target docker.service
Requires=docker.service

[Service]
User={{ fire2_user }}
Group={{ fire2_group }}
EnvironmentFile={{ fire2_install_dir }}/worker-%i.env
ExecStart={{ fire2_install_dir }}/worker
Restart=on-failure
RestartSec=5s

# Ensure the process group is killed on stop/restart — no orphan goroutines.
KillMode=mixed
TimeoutStopSec=10s

[Install]
WantedBy=multi-user.target
